# systemd/.config/systemd/user/look-out.service
[Unit]
Description=Reminder to look far and outside to reduce eyestrain, via libnotify

[Service]
Type=simple
ExecStart=/home/tiniyo/bin/lookout.zsh
