# systemd/.config/systemd/user/low-battery-notify.service
[Unit]
Description=Low battery notification via libnotify

[Service]
Type=simple
ExecStart=/home/prichaud/bin/batterynotify.zsh