# systemd/.config/systemd/user/overheat-notify.service
[Unit]
Description=CPU overheating notification via libnotify

[Service]
Type=simple
ExecStart=/home/tiniyo/bin/overheatnotify.zsh